# systemd/netbackup-backup.service
[Unit]
Description=NetBackup device configuration backup

[Service]
Type=oneshot
WorkingDirectory=/path/to/atrisDev/netbackup
ExecStart=/path/to/atrisDev/netbackup/scripts/run_backup.sh
