# systemd/netbackup-backup.service
[Unit]
Description=NetBackup device configuration backup

[Service]
Type=oneshot
WorkingDirectory=/path/to/netbackup
Environment=PYTHONPATH=src
ExecStart=/path/to/netbackup/scripts/run_backup.sh
